# integrations/opencode/services/opencode-web.service
[Unit]
Description=OpenCode Web Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/home/jon/.opencode/bin/opencode web
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target