# integrations/opencode/services/opencode-web.service
[Unit]
Description=OpenCode Web Server
After=network.target

Environment="OPENCODE_SERVER_PASSWORD=Mandielove2."

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=~/.opencode/bin/opencode web
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target