# scripts/service/templates/llama-server.service
[Unit]
Description=Llama.cpp Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__LLAMA_DIR__
ExecStart=__LLAMA_DIR__/build/bin/llama-server --host 127.0.0.1 --port 8080 --model __MODEL_PATH__
Restart=always
RestartSec=10
StandardOutput=append:/var/log/llama-server.log
StandardError=append:/var/log/llama-server.error.log

# Environment
Environment="GGML_CUDA_ENABLE_UNIFIED_MEMORY=1"

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
