# scripts/service/templates/llm-updater.service
[Unit]
Description=LLM Server Auto-Updater
After=network.target llm-frontend.service

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/scripts/update/monitor-updates.sh
Restart=always
RestartSec=60
StandardOutput=append:/var/log/llm-updater.log
StandardError=append:/var/log/llm-updater.error.log

[Install]
WantedBy=multi-user.target
