# services/llama.service
[Unit]
Description=OpenCode Web Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/home/jon/git/llm_server/scripts/launch_qwen3_35_b_4q.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target